# etc/systemd/system/postgresql.service
.include /usr/lib/systemd/system/postgresql.service
[Service]
Environment=PGDATA=/home/data/pgsql/data
